# postfix/etc/init/checkroot.sh.conf
# checkfs.sh - compatibility job for sysvinit dependencies
#
# This job runs once the root filesystem is mounted, to signal startpar
# that other rcS jobs relying on the historic checkroot.sh
# interface can continue.

description "Signal sysvinit that the rootfs is mounted"

start on mounted MOUNTPOINT=/
