# postfix/etc/init/dmesg.conf
# dmesg - save kernel messages
#
# This task saves the initial kernel message log.

description	"save kernel messages"

start on runlevel [2345]

task
script
    savelog -q -p -c 5 /var/log/dmesg
    dmesg -s 524288 > /var/log/dmesg
    chgrp adm /var/log/dmesg
end script
