# postfix/etc/init/plymouth.conf
# plymouth - Userspace bootsplash utility
#
# plymouth provides a boot splash screen on the system console using
# the kernel framebuffer device.  On boot, this is nominally started by
# the initramfs so the pre-start script, script and post-start script
# parts are actually not run.  These are normally run on shutdown instead.

description	"userspace bootsplash utility"

start on starting mountall

expect fork
kill timeout 60

exec /sbin/plymouthd --mode=boot --attach-to-session

pre-stop exec /bin/plymouth quit
