# postfix/etc/init/plymouth-splash.conf
# plymouth-splash - Show the splash screen
#
# plymouth must be started ASAP to avoid racing with gdm, but the splash
# screen can't be spawned until our framebuffer is available.  Wait for the
# video device to be available before showing the screen, or, if udevtrigger
# finishes without finding any video devices, bring up the fallback text
# interface.
# We also *should* wait for the filesystem to be up because of the libraries
# being used from /usr/lib, but this would cause a circular dependency if
# any interaction at all is required for mounting a filesystem; so these libs
# need to be moved to /lib instead.

description	"userspace bootsplash"

start on (started plymouth
          and (graphics-device-added PRIMARY_DEVICE_FOR_DISPLAY=1
               or drm-device-added PRIMARY_DEVICE_FOR_DISPLAY=1
               or stopped udev-fallback-graphics))

exec /bin/plymouth show-splash
